# docs/example_supervisor.ini
[program:serv1]
command=/usr/bin/node index.js
directory=/path/to/container-mgr-discord
autostart=true
autorestart=unexpected
user=user
redirect_stderr=true
environment=DISCORD_BOT_TOKEN="token",DOCKER_CONTAINER_ID="container ID"
